# deb-package/pkg/lib/systemd/system/dvoting.service
[Unit]
Description=D-Voting Node Service
After=network.target
Requires=unikernel.service

[Service]
User=dvoting
Type=forking
PIDFile=/var/dvoting.pid
ExecStart=/opt/dedis/dvoting/bin/start-dvoting
#ExecStart=echo "Starting fake DVoting"
ExecStop=/bin/kill --signal INT $MAINPID
#ExecStop=echo "Stopping fake DVoting"
#ExecReload=/bin/kill --signal HUP $MAINPID
KillMode=process
Restart=no
RestartSec=5

[Install]
WantedBy=multi-user.target
