# deb-package/pkg/lib/systemd/system/dvoting.service
[Unit]
Description=D-Voting Node Service
After=network.target
Requires=unikernel.service

[Service]
User=dvoting

ExecStartPre=/bin/bash -c "rm /var/opt/dedis/dvoting/data/dela/daemon.sock"
ExecStart=/bin/bash -c "/opt/dedis/dvoting/bin/start-dvoting /etc/dedis/dvoting/config.env"

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
