# packaging/pkg/lib/systemd/system/dvoting.service
[Unit]
Description=D-Voting Service
After=network.target

[Service]
User=dvoting
ExecStart=/opt/dedis/dvoting/bin/run /etc/dedis/dvoting/unikernel_paths
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=never

[Install]
WantedBy=multi-user.target
